# deployment/service.service
# systemd unit for the service. Update user and paths before installation.

[Unit]
Description=ServiceHandler
After=network.target

[Service]
Type=simple
User=yourusername
WorkingDirectory=/path/to/ServiceHandler
ExecStart=/path/to/ServiceHandler/.venv/bin/python src/main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target